# vpngateway/openvpn-uk-route-watchdog.timer
[Unit]
Description=Run OpenVPN UK remote route watchdog every 30 seconds

[Timer]
OnBootSec=2min
OnUnitActiveSec=30s
AccuracySec=5s
Persistent=true

[Install]
WantedBy=timers.target
