# homebridge-roborock/systemd/homebridge-network-watchdog.service
[Unit]
Description=LivingRoom network and Homebridge recovery watchdog
After=network-online.target homebridge.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/homebridge-network-watchdog
TimeoutStartSec=150s
